# initfiles/runit/ark_statistics/run
#!/bin/sh
exec chpst -u ark:ark -e /etc/sv/ark_statistics/env /home/ark/ark_xposed/ark_statistics.py

